# root@10.10.0.1111
#!/bin/bash
### BEGIN INIT INFO
# Provides:          init_node
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Inicia servidor
# Description:       Inicia servidor
### END INIT INFO

export HOME=/root
export MQTT_ADDR_LIGHTDNA=127.0.0.1
export MONGODB_ADDR_LIGHTDNA=127.0.0.1
export MONGODB_PORT_LIGHTDNA=27017
export KING_UNIVERSE=AIGNACIO
export LC_ALL=C

/opt/mongodb/bin/mongod --dbpath /dbpath_awges/ --repair --repairpath /dbpath_awges/
/opt/mongodb/bin/mongod --dbpath /dbpath_awges/ &
/usr/local/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf.example &
cd /root/lightdnav2
sudo /usr/local/bin/forever start server.js
sudo /usr/local/bin/forever start mqttPoll.js

